# agents/scheduler/crontab
# Tuner daily run: trigger at configured schedule, restart signals on success
# Default schedule: 00:00 UTC. Override via TUNER_CRON env var (per D-05).
0 0 * * * cd /app && docker compose -f docker-compose.yml --project-name phantom-perp run --rm tuner && docker compose -f docker-compose.yml --project-name phantom-perp restart signals >> /proc/1/fd/1 2>&1
# Required trailing newline for crond
